# contrib/ozwell-studio-launcher.service.in
[Unit]
Description=Ozwell Studio launcher — per-user container provisioning behind forward auth
Documentation=https://github.com/mieweb/ozwell-studio-launcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
EnvironmentFile=-/etc/default/ozwell-studio-launcher
Environment=NODE_ENV=production
WorkingDirectory=@PREFIX@
# Persistent pool state: $STATE_DIRECTORY (default SQL_URI puts db.sqlite here)
StateDirectory=ozwell-studio-launcher
ExecStart=/usr/bin/node @PREFIX@/src/server.js
Restart=on-failure
RestartSec=2

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK

[Install]
WantedBy=multi-user.target
